# backup/systemd/warden-verify.service
[Unit]
Description=warden verify (restic check)
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/warden verify
