# backup/systemd/warden-backup.service
[Unit]
Description=warden backup (restic)
After=docker.service network-online.target
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/warden backup
