# backup/systemd/warden-watch.service
[Unit]
Description=warden watch (monitoreo de contenedores)
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/warden watch
